# board/universal/rootfs-overlay/etc/systemd/system/oled-monitor.service
[Unit]
Description=OLED System Monitor
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/oled-monitor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
